# deploy/systemd/system/tracker-api.service
[Unit]
Description=OpenInfoMate API (Admin UI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/tracker
EnvironmentFile=/opt/tracker/.env
ExecStart=/opt/tracker/.venv/bin/tracker api serve
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
