# deploy/systemd/user/tracker.service
[Unit]
Description=OpenInfoMate (Info Mate)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/tracker
EnvironmentFile=%h/tracker/.env
ExecStart=%h/tracker/.venv/bin/tracker service run
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
